# var/lib/dpkg/info/ctxusb.conffiles
/opt/Citrix/ICAClient/usb.conf
/etc/init.d/ctxusbd
/etc/udev/rules.d/85-ica-usb.rules
