# var/lib/dpkg/info/cups-daemon.conffiles
/etc/ufw/applications.d/cups
/etc/pam.d/cups
/etc/logrotate.d/cups-daemon
/etc/apparmor.d/usr.sbin.cupsd
/etc/cups/cups-files.conf
/etc/init.d/cups
